# websocket-visor.conf
[program:websocket-visor]
command=/root/projects/OpenWorship/websocket.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
